# data/io.elementary.appcenter.service.in
[Unit]
Description=Check for package updates, and notify that updates are available

[Service]
ExecStart=/usr/lib/@PROJECT_NAME@/update-check.sh